# supervisord.conf
[supervisord]
nodaemon=true

[program:dnsmasq]
command=/usr/sbin/dnsmasq --no-daemon

[program:hostapd]
command=/usr/sbin/hostapd /etc/hostapd/hostapd.conf